# dist/lib/systemd/system/cloudjournal.service
[Unit]
Description=cloudjournal service

[Service]
User=root
Group=root
ExecStart=/usr/bin/cloudjournal start --config=/etc/cloudjournal/config.yaml
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
