# dist/DEBIAN/postinst
#!/bin/sh
systemctl enable cloudjournal.service
systemctl start cloudjournal.service
