# dist/DEBIAN/prerm
#!/bin/sh
systemctl stop cloudjournal.service
systemctl disable cloudjournal.service
